# automation/cron/generated.crontab
# aigraph 24/7 automation schedule
# repo: /workspace/literature-conflict-graph
# generated: 2026-04-17T17:23:44
SHELL=/bin/bash
PATH=/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin

15 * * * * cd /workspace/literature-conflict-graph && /workspace/literature-conflict-graph/.venv/bin/python -m aigraph.cli automation-harvest --automation-dir /workspace/literature-conflict-graph/automation --runs-dir /workspace/literature-conflict-graph/outputs/runs >> /workspace/literature-conflict-graph/automation/logs/topic_harvest.log 2>&1
25 * * * * cd /workspace/literature-conflict-graph && /workspace/literature-conflict-graph/.venv/bin/python -m aigraph.cli automation-run-batch --automation-dir /workspace/literature-conflict-graph/automation --runs-dir /workspace/literature-conflict-graph/outputs/runs --batch-size 3 >> /workspace/literature-conflict-graph/automation/logs/run_batch.log 2>&1
40 */2 * * * cd /workspace/literature-conflict-graph && /workspace/literature-conflict-graph/.venv/bin/python -m aigraph.cli automation-critic --automation-dir /workspace/literature-conflict-graph/automation --runs-dir /workspace/literature-conflict-graph/outputs/runs --limit 8 >> /workspace/literature-conflict-graph/automation/logs/critic_batch.log 2>&1
50 2 * * * cd /workspace/literature-conflict-graph && /workspace/literature-conflict-graph/.venv/bin/python -m aigraph.cli automation-fix-bundle --automation-dir /workspace/literature-conflict-graph/automation --runs-dir /workspace/literature-conflict-graph/outputs/runs --max-issues 3 >> /workspace/literature-conflict-graph/automation/logs/fix_bundle.log 2>&1
10 3 * * * cd /workspace/literature-conflict-graph && if [ -n "$AIGRAPH_CODEX_FIX_COMMAND" ]; then /workspace/literature-conflict-graph/.venv/bin/python -m aigraph.cli automation-fix-run --automation-dir /workspace/literature-conflict-graph/automation --runs-dir /workspace/literature-conflict-graph/outputs/runs --repo-dir /workspace/literature-conflict-graph --max-issues 3 --push --open-pr; else echo 'Skipping automation-fix-run: AIGRAPH_CODEX_FIX_COMMAND is not set'; fi >> /workspace/literature-conflict-graph/automation/logs/fix_run.log 2>&1
